# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/data/logs/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/data/logs/nginx_stdout.log
stderr_logfile=/data/logs/nginx_stderr.log

[program:solar-assistant]
command=python /app/solar_assistant_server.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/data/logs/solar_assistant_stdout.log
stderr_logfile=/data/logs/solar_assistant_stderr.log
environment=PYTHONUNBUFFERED="1"

[program:mqtt-bridge]
command=python /app/mqtt_bridge.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/data/logs/mqtt_bridge_stdout.log
stderr_logfile=/data/logs/mqtt_bridge_stderr.log

[program:data-collector]
command=python /app/data_collector.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/data/logs/data_collector_stdout.log
stderr_logfile=/data/logs/data_collector_stderr.log

[group:solar]
programs=nginx,solar-assistant,mqtt-bridge,data-collector